# deployments/puppet/templates/ld.so.preload.erb
/usr/lib/splunk-instrumentation/libsplunk.so
<% if defined?(@auto_instrumentation_ld_so_preload) && @auto_instrumentation_ld_so_preload != "" -%>
<%= @auto_instrumentation_ld_so_preload %>
<% end -%>
