# deployments/puppet/templates/ld.so.preload.erb
<% if defined?(@auto_instrumentation_systemd) && @auto_instrumentation_systemd.to_s.downcase == "false" -%>
/usr/lib/splunk-instrumentation/libsplunk.so
<% end -%>
<% if defined?(@auto_instrumentation_ld_so_preload) && @auto_instrumentation_ld_so_preload != "" -%>
<%= @auto_instrumentation_ld_so_preload %>
<% end -%>
