# deployments/ansible_collections/signalfx/splunk_otel_collector/roles/gateway/templates/gateway.service.j2
# OpenTelemetry Collector Gateway Service Override
# Generated by Ansible - DO NOT EDIT MANUALLY

[Service]
# Gateway mode requires higher resource limits
EnvironmentFile={{ otel_gateway_env_file }}
User={{ otel_gateway_user }}
Group={{ otel_gateway_group }}

# Resource limits for gateway deployment
LimitNOFILE=8192
LimitNPROC=4096

# Gateway-specific command line
ExecStart=
ExecStart=/usr/bin/otelcol --config=${SPLUNK_CONFIG}

# Restart policy
Restart=on-failure
RestartSec=5s
