# deployments/chef/templates/ld.so.preload.erb
<% if !@with_systemd && (@with_java || @with_nodejs) -%>
/usr/lib/splunk-instrumentation/libsplunk.so
<% end -%>
<% if defined?(node['splunk_otel_collector']['auto_instrumentation_ld_so_preload']) && node['splunk_otel_collector']['auto_instrumentation_ld_so_preload'] != "" -%>
<%= node['splunk_otel_collector']['auto_instrumentation_ld_so_preload'] %>
<% end -%>
